# claude-plan-api.service
[Unit]
Description=Claude Plan API — OpenAI-compatible gateway for Claude subscriptions
After=network.target

[Service]
Type=simple
ExecStart=/root/.bun/bin/bun run /opt/claude-plan-api/src/index.ts 3456
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
